# service/goMarkableStream.service
[Unit]
Description=goMarkableStream

[Service]
ExecStart=/home/root/goMarkableStreamServer.arm
Restart=always

[Install]
WantedBy=multi-user.target